# files/etc/supervisor.d/openvpn.ini
[program:openvpn]
command = /usr/sbin/openvpn --cd /vpn --config client.conf
